# package/external/pdnsd/files/pdnsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) & wongsyrone

. /lib/functions.sh

START=95

EXECUTEBIN="usr/sbin/pdnsd"
PIDFILE="/var/run/pdnsd.pid"
INIFILE="/etc/pdnsd.conf"
CACHEDIR=/var/pdnsd
CACHE=$CACHEDIR/pdnsd.cache
ENABLED=0
USER=nobody
GROUP=nogroup

config_load "pdnsd"

parse_pdnsd()
{
	config_get ENABLED "$section" "enabled"
}
config_foreach parse_pdnsd 'pdnsd'

gen_cache() {
	[ -f "$CACHE" ] || {
		mkdir -p $CACHEDIR
		dd if=/dev/zero of="$CACHE" bs=1 count=4 2> /dev/null
		chown -R $USER.$GROUP $CACHEDIR
	}
}

start() {
	if [ "1" = "$ENABLED" ] || [ "on" = "$ENABLED" ] || [ "true" = "$ENABLED" ]; then
		gen_cache
		$EXECUTEBIN --daemon -p $PIDFILE
	else
		echo "pdnsd is disabled!"
	fi
}

stop() {
	[ -f $PIDFILE ] && {
		kill `cat $PIDFILE`
		echo "pdnsd terminated successfully"
		rm -f $PIDFILE
	} || {
		echo "pdnsd is not running"
	}
}

restart() {
	stop
	sleep 2;
	start
}
