# package/external/transocks-wong/files/transocks-wong.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2018 wongsyrone

. /lib/functions.sh

START=95
USE_PROCD=1
#PROCD_DEBUG=1


PROG=/usr/sbin/transocks-wong

config_load "transocks-wong"

parse_transocks_wong()
{
	config_get ENABLED "$section" "enabled" "0"
	config_get LISTENERADDRPORT "$section" "listener_addr_port" "[::]:8123"
	config_get SOCKS5ADDRPORT "$section" "socks5_addr_port" "[::1]:1081"
	config_get PUMPMETHOD "$section" "pump_method" "splicepump"
}

config_foreach parse_transocks_wong 'transocks-wong'


start_service() {
	if [ "1" = "$ENABLED" ] || [ "on" = "$ENABLED" ] || [ "true" = "$ENABLED" ]; then
		local cpunum=$(grep ^cpu\\scores /proc/cpuinfo | uniq |  awk '{print $4}' )
		for i in $(seq 1 ${cpunum});
		do
			procd_open_instance
			procd_set_param command $PROG
			procd_append_param command --listener-addr-port="$LISTENERADDRPORT"
			procd_append_param command --socks5-addr-port="$SOCKS5ADDRPORT"
			procd_append_param command --pump-method="$PUMPMETHOD"
			procd_set_param user root # run service as user root
			procd_set_param stdout 1 # forward stdout of the command to logd
			procd_set_param stderr 1 # same for stderr
			procd_set_param limits nofile="unlimited"
			[ -e /proc/sys/kernel/core_pattern ] && {
				procd_append_param limits core="unlimited"
			}
			procd_close_instance
		done
	else
		echo "transocks-wong is disabled"
	fi
}
