# package/external/transocks-wong/files/transocks-wong.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2018 wongsyrone

. /lib/functions.sh

START=95
USE_PROCD=1
#PROCD_DEBUG=1


PROG=/usr/sbin/transocks-wong
LISTENARG=--listener-addr-port="[::]:8123"
TARGETARG=--socks5-addr-port="[::1]:1081"
PUMPARG=--pump-method=splicepump

config_load "transocks-wong"

parse_transocks_wong()
{
	config_get ENABLED "$section" "enabled"
}

config_foreach parse_transocks_wong 'transocks-wong'


start_service() {
	if [ "1" = "$ENABLED" ] || [ "on" = "$ENABLED" ] || [ "true" = "$ENABLED" ]; then
		procd_open_instance
		procd_set_param command $PROG
		procd_append_param command $LISTENARG
		procd_append_param command $TARGETARG
		#procd_append_param command $PUMPARG
		procd_set_param user root # run service as user root
		procd_set_param stdout 1 # forward stdout of the command to logd
		procd_set_param stderr 1 # same for stderr
		procd_set_param limits nofile="unlimited"
		[ -e /proc/sys/kernel/core_pattern ] && {
			procd_append_param limits core="unlimited"
		}
		procd_close_instance
	else
		echo "transocks-wong is disabled"
	fi
}
