# package/external/luci-app-zerotier/root/etc/hotplug.d/net/99-zerotier
#!/bin/sh

#set -x

. "$IPKG_INSTROOT"/lib/functions.sh
. /etc/init.d/zerotier

NFTABLES_AUTO_NAT_RULE_FILE="/etc/nftables.d/99-zerotier-auto-nat-auto-generated-do-not-touch.nft"

my_log() {
	logger -t zerotier-hotplug "$@" "ACTION=$ACTION DEVPATH=$DEVPATH SUBSYSTEM=$SUBSYSTEM INTERFACE=$INTERFACE"
}

# $1: timeout
delay_run_in_background() {
	local timeout="$1"
	shift
	sleep $timeout;
	eval "$*" &
}

add_iptables_rule() {
	local zt_iface="$1"
	my_log "add iptables rule"
	iptables -I FORWARD -i $zt_iface -j ACCEPT -m comment --comment "zt_hotplug: allow inbound forward $zt_iface"
	iptables -I FORWARD -o $zt_iface -j ACCEPT -m comment --comment "zt_hotplug: allow outbound forward $zt_iface"
	iptables -t nat -I POSTROUTING -o $zt_iface -j MASQUERADE -m comment --comment "zt_hotplug: allow outbound postrouting masq $zt_iface"
}

del_iptables_rule() {
	local zt_iface="$1"
	my_log "del iptables rule"
	iptables -D FORWARD -i $zt_iface -j ACCEPT
	iptables -D FORWARD -o $zt_iface -j ACCEPT
	iptables -t nat -D POSTROUTING -o $zt_iface -j MASQUERADE
}

add_nftables_rule() {
	local zt_iface="$1"
	cat > $NFTABLES_AUTO_NAT_RULE_FILE <<EOF
	chain zt_auto_nat_pre_fwd {
		type filter hook forward priority -50; policy accept;
		iifname "$zt_iface" counter accept comment "zt_hotplug: allow inbound forward $zt_iface"
		oifname "$zt_iface" counter accept comment "zt_hotplug: allow outbound forward $zt_iface"
	}
	chain zt_auto_nat_post_rt {
		type nat hook postrouting priority srcnat; policy accept;
		oifname "$zt_iface" counter meta nfproto ipv6 masquerade comment "zt_hotplug: allow outbound postrouting masq IPv6 $zt_iface"
		oifname "$zt_iface" counter meta nfproto ipv4 masquerade comment "zt_hotplug: allow outbound postrouting masq IPv4 $zt_iface"
	}
EOF
	my_log "add nftables rule"
	/etc/init.d/firewall reload
}

del_nftables_rule() {
	local zt_iface="$1"
	rm -f $NFTABLES_AUTO_NAT_RULE_FILE
	nft flush chain inet fw4 zt_auto_nat_pre_fwd
	nft flush chain inet fw4 zt_auto_nat_post_rt
	nft delete chain inet fw4 zt_auto_nat_pre_fwd
	nft delete chain inet fw4 zt_auto_nat_post_rt
	my_log "del nftables rule"
	/etc/init.d/firewall reload
}

if [ "$SUBSYSTEM" != "net" ]; then
	my_log "not net event, ignore"
	exit 0
fi

case $INTERFACE in
	"zt"*)
	my_log "found zerotier interface: $INTERFACE"
	;;
	*)
	# not found
	my_log "not zerotier interface: $INTERFACE, ignore"
	exit 0
	;;
esac

# adding the unmanaged interface deletes the zerotier entry from the routing table.
# leave & join fixes it
# 'resyncnetwork' works with patched zerotier
resyncnetwork() {
	local netId="$1"
	local ret
	ret=$(zerotier-cli resyncnetwork "$netId")
	my_log "sync network managed routes: $netId [$ret]"
}

set_network_local_conf_cli() {
	local netId="$1"
	local paramName="$2"
	local paramValue="$3"
	local ret
	zerotier-cli set "$netId" "$paramName"="$paramValue" > /dev/null
	# checks
	ret=$(zerotier-cli get "$netId" "$paramName")
	if [ x"$ret" = x"$paramValue" ]; then
		my_log "set $paramName to $paramValue successfully"
	else
		my_log "set $paramName to $paramValue fail"
	fi
}


openwrt_network_auto_nat_enabled() {
	nat_enable="$(uci get zerotier.openwrt_network.nat)"
	[ "$nat_enable" = "1" ]
}

remove_netif() {
	uci -q batch <<-EOF >/dev/null
delete network.$INTERFACE
commit network
EOF
	my_log "remove interface"
}

remove_firewall_rules() {
	my_log "delete firewall rules"
	if [ -f "/usr/sbin/nft" ]; then
		del_nftables_rule "$INTERFACE"
	else
		del_iptables_rule "$INTERFACE"
	fi

	uci -q batch <<-EOF >/dev/null
delete firewall.zt_zone_$INTERFACE
delete firewall.zt_rule_allowInbound
delete firewall.zt_fwd_${INTERFACE}2lan
delete firewall.zt_fwd_${INTERFACE}2wan
delete firewall.zt_fwd_lan2${INTERFACE}
delete firewall.zt_rule_allowSecondPort
delete firewall.zt_rule_allowThirdPort
commit firewall
EOF
}

if [ "$ACTION" = "remove" ]; then
	remove_firewall_rules
	remove_netif
	/etc/init.d/network reload
fi

if [ "$ACTION" = "add" ]; then
		remove_firewall_rules

		uci -q batch <<-EOF
set network.$INTERFACE=interface
set network.$INTERFACE.ifname=$INTERFACE
set network.$INTERFACE.proto=none
set network.$INTERFACE.device=$INTERFACE
commit network
add firewall zone
rename firewall.@zone[-1]=zt_zone_$INTERFACE
set firewall.@zone[-1].network=$INTERFACE
set firewall.@zone[-1].input=ACCEPT
set firewall.@zone[-1].output=ACCEPT
set firewall.@zone[-1].forward=ACCEPT
set firewall.@zone[-1].name=zerotier
set firewall.@zone[-1].masq=1
set firewall.@zone[-1].mtu_fix=0
add_list firewall.@zone[-1].device=$INTERFACE
commit firewall
EOF
		config_load 'zerotier'
		
		netIds="$(uci get zerotier.openwrt_network.join)"
		for netId in ${netIds}; do
			myVal=
			for pName in allowDNS allowDefault allowGlobal; do
				config_get_bool myVal "openwrt_network" "$pName" 0
				set_network_local_conf_cli "$netId" "$pName" "$myVal"
			done
			
			config_get_bool myVal "openwrt_network" "allowManaged" 1
			set_network_local_conf_cli "$netId" "allowManaged" "$myVal"

			delay_run_in_background 15 resyncnetwork "$netId"
		done
		
		if openwrt_network_auto_nat_enabled; then
			my_log "auto NAT enabled; add firewall config"
			
			port="$(uci get zerotier.openwrt_network.port)"
			[ -z "$port" ] && port=9993
			firstPort=$(cat "$FIRST_PORT_FILE")
			secondPort=$(cat "$SECOND_PORT_FILE")
			thirdPort=$(cat "$THIRD_PORT_FILE")
			if [ -n "$port" ] && [ -n "$firstPort" ] && [ x"$port" != x"$firstPort" ]; then
				my_log "logic error, port different $port and $firstPort"
			fi
			
			uci -q batch <<-EOF
add firewall rule
rename firewall.@rule[-1]=zt_rule_allowInbound
set firewall.@rule[-1].name='Allow-ZeroTier-Inbound'
set firewall.@rule[-1].src='*'
set firewall.@rule[-1].target='ACCEPT'
add_list firewall.@rule[-1].proto='all'
set firewall.@rule[-1].dest_port='$port'
commit firewall
add firewall forwarding
rename firewall.@forwarding[-1]=zt_fwd_${INTERFACE}2lan
set firewall.@forwarding[-1].src=zerotier
set firewall.@forwarding[-1].dest=lan
add firewall forwarding
rename firewall.@forwarding[-1]=zt_fwd_${INTERFACE}2wan
set firewall.@forwarding[-1].src=zerotier
set firewall.@forwarding[-1].dest=wan
add firewall forwarding
rename firewall.@forwarding[-1]=zt_fwd_lan2${INTERFACE}
set firewall.@forwarding[-1].src=lan
set firewall.@forwarding[-1].dest=zerotier
commit firewall
EOF
			[ -n "$secondPort" ] && {
			uci -q batch <<-EOF
add firewall rule
rename firewall.@rule[-1]=zt_rule_allowSecondPort
set firewall.@rule[-1].name='Allow-ZeroTier-2ndPort'
set firewall.@rule[-1].src='*'
set firewall.@rule[-1].target='ACCEPT'
add_list firewall.@rule[-1].proto='all'
set firewall.@rule[-1].dest_port='$secondPort'
commit firewall
EOF
			}

			[ -n "$thirdPort" ] && {
			uci -q batch <<-EOF
add firewall rule
rename firewall.@rule[-1]=zt_rule_allowThirdPort
set firewall.@rule[-1].name='Allow-ZeroTier-3rdPort'
set firewall.@rule[-1].src='*'
set firewall.@rule[-1].target='ACCEPT'
add_list firewall.@rule[-1].proto='all'
set firewall.@rule[-1].dest_port='$thirdPort'
commit firewall
EOF
			}
			if [ -f "/usr/sbin/nft" ]; then
				add_nftables_rule "$INTERFACE"
			else
				add_iptables_rule "$INTERFACE"
			fi
		fi
fi
